# v-todo/tmp-archives/dot/profile.d/omz.sh
# oh-my-zsh
export ZSH="$HOME/Applications/.oh-my-zsh"

ZSH_THEME="robbyrussell"

source $ZSH/oh-my-zsh.sh

plugins=(
    z
    zsh-autosuggestions
    git-open
)

# alias zshconfig="mate ~/.zshrc"
# alias ohmyzsh="mate ~/.oh-my-zsh"
