# v-todo/tmp-archives/dot/profile.d/zconf.sh
# 放一些命令设置

# history 显示时间
HIST_STAMPS="yyyy-mm-dd"
